# ual.crontab
56	*/4	*	*	*	ual/src/ual.py ual/alerts/alert_defs.txt 1>> ual/logs/ual_search.log 2>> ual/logs/ual_search.err
45	8,20	*	*	*	ual/src/ual.py ual/alerts/agg_alert_defs.txt -a 1>> ual/logs/ual_long_search.log 2>> ual/logs/ual_long_search.err
23	*	*	*	*	ual/src/ual.py ual/alerts/text_alerts.txt -t 1>> ual/logs/ual_text_search.log 2>> ual/logs/ual_text_search.err
*/4	*	6-7	11	*	ual/src/ual.py ual/alerts/fast_alerts.txt -t 1>> ual/logs/ual_fast_search.log 2>> ual/logs/ual_fast_search.err
36	*/2	*	*	*	ual/src/ual.py ual/alerts/nan_alerts.txt -e 6507961064@txt.att.net 1>> ual/logs/nan_search.log 2>> ual/logs/nan_search.err
